# templates/gitlab.backup.cron.j2
# m  h          d M W
  30 4          * * *           root    /usr/bin/docker exec -t {{container_name}} gitlab-rake gitlab:backup:create
